# ltsp-org-updater/etc/systemd/system/ltsp-org-update.service
[Unit]
Description=Update Images

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ltsp-org-updater --update

[Install]
WantedBy=multi-user.target
