# ltsp-org-updater/etc/systemd/system/ltsp-org-updater.service
[Unit]
Description=Update Images

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ltsp-org-updater --update

[Install]
WantedBy=multi-user.target
