# ltsp-org-updater/etc/systemd/system/ltsp-org-initrd.service
[Unit]
Description=Configure initrd

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ltsp-org-updater --initrd

[Install]
WantedBy=multi-user.target
