# ltsp-org-updater/etc/systemd/system/ltsp-org-daemons.service
[Unit]
Description=Update Daemons

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ltsp-org-updater --daemons

[Install]
WantedBy=multi-user.target
