# ltsp-org-updater/etc/systemd/system/ltsp-org-config.service
[Unit]
Description=Configure LTSP-Server

[Service]
ExecStart=/usr/local/bin/ltsp-org-updater --daemon

[Install]
WantedBy=multi-user.target