# ltsp-updater/etc/systemd/system/ltsp-org-update.service
[Unit]
Description=Update Image

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ltsp-updater

[Install]
WantedBy=multi-user.target